# settings/autopilot-board.service.template
# Autopilot Board — systemd --user unit (managed by airuleset).
#
# Installed ONLY on the board host (dev1, 10.77.9.21) by `airuleset.py install`
# via setup_board_service(). The {{REPO_DIR}} placeholder is substituted with the
# absolute path of THIS airuleset checkout at install time, so ExecStart always
# points at the right airuleset.py.
#
# Do not edit the installed copy in ~/.config/systemd/user/ by hand — edit this
# template and re-run `python3 airuleset.py install`.

[Unit]
Description=Autopilot Board — central live tracking + review-gate audit
Documentation=https://github.com/zbynekdrlik/airuleset
After=network-online.target
Wants=network-online.target
# Restart loop guard: more than StartLimitBurst restarts in the interval = stop.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/env python3 {{REPO_DIR}}/airuleset.py board --serve
Restart=on-failure
RestartSec=5

# --- Hardening (defense in depth; the board is LAN-internal, token-gated) ---
NoNewPrivileges=yes
ProtectSystem=strict
# The board writes its sqlite DB, token, and queue under ~/.claude.
ReadWritePaths=%h/.claude
ProtectHome=read-only
PrivateTmp=yes
# Network sockets only — the board speaks HTTP over TCP and shells out to gh.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
MemoryMax=512M
TasksMax=64

[Install]
WantedBy=default.target
