# settings/filedrop.service.template
# airuleset File-Drop — systemd --user unit (managed by airuleset).
#
# Installed on EVERY machine by `airuleset.py install` via setup_filedrop_service()
# (unlike the autopilot board, which is dev1-only). Each machine serves the files
# produced on THAT machine, bound to THAT machine's own LAN IP. The {{REPO_DIR}}
# placeholder is substituted with the absolute path of THIS airuleset checkout at
# install time, so ExecStart always points at the right airuleset.py.
#
# Do not edit the installed copy in ~/.config/systemd/user/ by hand — edit this
# template and re-run `python3 airuleset.py install`.

[Unit]
Description=airuleset File-Drop — serve user files as clickable LAN URLs
Documentation=https://github.com/zbynekdrlik/airuleset
After=network-online.target
Wants=network-online.target
# Restart loop guard: more than StartLimitBurst restarts in the interval = stop.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# The IPs to bind are baked in at install time (computed UNSANDBOXED, where
# `hostname -I` / `tailscale ip` work). The sandbox below restricts address families
# to AF_INET/AF_INET6, which blocks the AF_NETLINK socket `hostname -I` needs — so
# the server must NOT discover its own IP at runtime.
#   FILEDROP_HOSTS = comma list of ALL private interfaces (tailscale + LAN) the
#     server binds, so the file-drop answers whether the user is on tailscale or the
#     LAN. NEVER the box's public IP (bind_ips() excludes it). A bind that fails
#     (stale LAN IP after a network switch) is skipped, not fatal — tailscale stays up.
#   FILEDROP_HOST = the primary (tailscale-first) IP, kept for status/URL + as the
#     single-host fallback if FILEDROP_HOSTS is empty.
# Both are substituted per machine.
Environment=FILEDROP_HOSTS={{HOST_IPS}}
Environment=FILEDROP_HOST={{HOST_IP}}
# Per-user port: a second airuleset user on the same host (montalu@dev1,
# marek@gatekeeper) cannot reuse the first user's 8788 — install picks a free
# port and persists it to ~/.claude/filedrop.port so the share CLI agrees.
Environment=FILEDROP_PORT={{PORT}}
ExecStart=/usr/bin/env python3 {{REPO_DIR}}/airuleset.py filedrop --serve
Restart=on-failure
RestartSec=5

# --- Hardening (defense in depth; LAN-internal, read-only file server) ---
NoNewPrivileges=yes
ProtectSystem=strict
# Read-only server: it only READS ~/.claude/filedrop to serve files. The `share`
# client (runs unsandboxed as the user) is what writes/prunes there, so the
# service needs NO write paths — home is readable but not writable.
ProtectHome=read-only
PrivateTmp=yes
# TCP sockets only — plain static HTTP, no subprocesses, no unix sockets.
RestrictAddressFamilies=AF_INET AF_INET6
MemoryMax=256M
TasksMax=64

[Install]
WantedBy=default.target
